# packaging/linux/morecode.service
[Unit]
Description=more Code headless server
Documentation=https://github.com/pingdotgg/t3code
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=morecode
Group=morecode
WorkingDirectory=/var/lib/morecode/workspace
Environment=HOME=/var/lib/morecode
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
EnvironmentFile=-/etc/morecode/morecode.env
ExecStart=/usr/bin/morecode serve
Restart=on-failure
RestartSec=5s
TimeoutStopSec=30s
KillSignal=SIGTERM
UMask=0027
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full

[Install]
WantedBy=multi-user.target

